# Config/Seriot/etc/systemd/system/snap-gimp-393.mount
[Unit]
Description=Mount unit for gimp, revision 393
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/gimp_393.snap
Where=/snap/gimp/393
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
