# Config/Seriot6/etc/systemd/system/multi-user.target.wants/snap-gnome\x2dlogs-112.mount
[Unit]
Description=Mount unit for gnome-logs, revision 112
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/gnome-logs_112.snap
Where=/snap/gnome-logs/112
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
