# Config/Seriot5/etc/systemd/system/timers.target.wants/portscan.timer
[Unit]
Description=Perform port scan every 30+1 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=31min

[Install]
WantedBy=timers.target
