# Config/Seriot4/etc/network/if-up.d/br0
#!/bin/sh
ip addr add dev br0 10.0.1.4/24
ip l set up dev br0
arp -f

