# Config/Seriot5/etc/systemd/system/portscan.service
[Unit]
Description=Perform portscan

[Service]
Type=oneshot
ExecStart=/usr/bin/nmap -vvv -p 1-1000 10.0.2.33
